# docker/supervisord.dev.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── FastAPI (uvicorn with hot-reload) ─────────────────────────────────────────
; WATCHFILES_FORCE_POLLING: inotify doesn't fire on Windows+Docker bind mounts;
; polling ensures uvicorn --reload detects saves from the Windows host.
[program:backend]
command=/app/backend/.venv/bin/uvicorn app.main:app
    --reload
    --host 127.0.0.1
    --port 8000
environment=WATCHFILES_FORCE_POLLING="true",WATCHFILES_POLL_DELAY_MS="500"
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Next.js dev server (webpack + polling) ────────────────────────────────────
; --no-turbopack: Turbopack has no polling mode — on Windows+Docker bind mounts
; it never sees file changes. Webpack respects WATCHPACK_POLLING.
; CHOKIDAR_USEPOLLING: fallback for any chokidar-based watcher inside Next.js.
[program:frontend]
command=bun run next dev --webpack --hostname 127.0.0.1 --port 3000
environment=WATCHPACK_POLLING="true",CHOKIDAR_USEPOLLING="true",CHOKIDAR_INTERVAL="500"
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Nginx reverse proxy ────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
