# docker/supervisord.prod.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── FastAPI (Gunicorn + Uvicorn workers) ─────────────────────────────────────
[program:backend]
command=/app/backend/.venv/bin/gunicorn app.main:app
    -k uvicorn.workers.UvicornWorker
    --bind 127.0.0.1:8000
    --workers 2
    --timeout 120
    --access-logfile -
    --error-logfile -
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Next.js standalone server ─────────────────────────────────────────────────
[program:frontend]
command=node /app/frontend/server.js
environment=PORT=3000,HOSTNAME="127.0.0.1",NODE_ENV="production"
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Nginx reverse proxy ───────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
